# task_local_15min.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>DESKTOP-AUSHUPM\ADMIN</Author>
    <URI>\LocalBackup_15min</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <StartWhenAvailable>true</StartWhenAvailable>
    <Enabled>true</Enabled>
  </Settings>
  <Triggers>
    <TimeTrigger>
      <StartBoundary>2026-04-12T00:00:00</StartBoundary>
      <Repetition>
        <Interval>PT15M</Interval>
      </Repetition>
    </TimeTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c "d:\0 - Google Antigravity\11 - NHAN VIEN KINH DOANH - Copy\backup_local.bat"</Arguments>
      <WorkingDirectory>d:\0 - Google Antigravity\11 - NHAN VIEN KINH DOANH - Copy</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
